# etc/systemd/system/docker.service.d/flannel.conf
[Unit]
After=flannel.service
Requires=flannel.service

[Service]
EnvironmentFile=/run/flannel/subnet.env
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --mtu=${FLANNEL_MTU} --bip=${FLANNEL_SUBNET}
